# xinetd.d/xinetd.d/daytime-lo
service daytime-lo
{
  disable          = no
  type             = UNLISTED
  socket_type      = stream
  protocol         = tcp
  wait             = no
  redirect         = localhost 13
  bind             = 127.0.0.1
  port             = 4013
  user             = nobody
}

